# configs/services/timevox-shutdown.service
[Unit]
Description=TimeVox Shutdown Button Service
After=multi-user.target
Conflicts=halt.target reboot.target shutdown.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/timevox/shutdown_button.py
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
